# eg/simple.crontab
MAILTO=test@example.com

* * * * * echo "HELLO"

