# eg/simple.crontab
MAILTO=test@example.com
FOO=bar

* * * * * echo "HELLO"

